# debian/postinst
#!/bin/sh

systemctl enable rockchip-chromium-x11-init